# testing/aws-hostname/aws-hostname.initd
#!/sbin/openrc-run

NAME=aws-hostname
DAEMON=/usr/bin/$NAME

depend() {
        need net
}

start() {
        ebegin "Starting ${NAME}"
                $DAEMON
                /sbin/rc-update del ${NAME}
        eend $?
}