# testing/headstart/headstart.initd
#!/sbin/openrc-run

depend() {
	need net
	provide headstart
	use dns
}

command="/usr/bin/headstart"